# share/check_argc0_lsm.service.tmpl
[Unit]
Description=LSM BPF protection for argc == 0 execs

[Service]
Type=simple
Restart=on-failure
ExecStart=DESTDIR/sbin/load_check_argc0_lsm

[Install]
WantedBy=multi-user.target
